# scripts/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ── MariaDB ──────────────────────────────────────────────────────────────
[program:mariadb]
command=/usr/bin/mysqld_safe --datadir=/data/mysql
user=root
autostart=true
autorestart=true
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startretries=3

; ── Flask / Gunicorn ─────────────────────────────────────────────────────
[program:flask]
command=gunicorn web.app:app --bind 0.0.0.0:%(ENV_PORT)s --workers 2 --timeout 120
directory=/app
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/root",PATH="/usr/local/bin:%(ENV_PATH)s"
startretries=5
startsecs=8

; ── FastAPI / Uvicorn ────────────────────────────────────────────────────
[program:fastapi]
command=uvicorn api.main:app --host 0.0.0.0 --port %(ENV_API_PORT)s --log-level warning
directory=/app
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/root",PATH="/usr/local/bin:%(ENV_PATH)s"
startretries=5
startsecs=8
